# highwaydragons.conf
[program:highwaydragons]
command = /srv/highwaydragons/gunicorn_start.bash
user = root
stdout_logfile = /srv/highwaydragons/log/gunicorn_stdout.log
redirect_stderr = true
