# highwaydragons.conf
[program:highwaydragons]
command = /var/www/HighwayDragons/gunicorn_start.bash
user = root
stdout_logfile = /var/www/HighwayDragons/log/gunicorn_stdout.log
redirect_stderr = true
